# supervisor/apiok-admin.ini
[program:apiok-admin]
command=/opt/apiok-admin/apiok-admin
directory=/opt/apiok-admin
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/opt/apiok-admin/logs/stdout.log
stderr_logfile=/opt/apiok-admin/logs/stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=GIN_MODE=release

